# Library/LaunchAgents/local.appearance-change.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>local.appearance-change</string>
    <key>RunAtLoad</key>
    <true/>
    <key>WatchPaths</key>
    <array>
      <string>~/Library/Preferences/.GlobalPreferences.plist</string>
    </array>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/env</string>
      <string>bash</string>
      <string>-c</string>
      <string>source ~/.config/dotfiles/appearance-change/appearance-change.bash</string>
    </array>
  </dict>
</plist>
